# .repofiles/compiled-packages/apt.1nternist.bashrc_1.0_iphoneos-arm/private/etc/profile.d/profile.sh
if [ "${BASH-no}" != "no" ]; then
    [ -r /etc/bashrc ] && . /etc/bashrc
fi
